# templates/kafka-exporter.service.j2
[Unit]
Description=Kafka Exporter for Prometheus
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ user }}
Group={{ group }}
ExecReload=/bin/kill -HUP
ExecStart=/usr/local/bin/kafka_exporter \
	{% for broker in kafka_brokers -%}
	--kafka.server={{ broker }} \
	{% endfor %}
	--web.listen-address=:{{ exporter_port }}
SyslogIdentifier=kafka-exporter
Restart=always

[Install]
WantedBy=multi-user.target
